# mail/.config/systemd/user/offlineimap-oneshot.service.d/override.conf
[Service]
ExecStart=
ExecStart=offlineimap -o -u basic -a Ens,Proton,Inria