# workflow/.config/systemd/user/battery-low.service
[Unit]
Description=Issue a battery low alert if battery is under 10%

[Service]
Type=simple
ExecStart=/home/jean/.local/bin/battery_low_alert.sh
