# mail/.config/systemd/user/offlineimap-oneshot.timer.d/override.conf
[Timer]
OnUnitInactiveSec=
OnUnitInactiveSec=30m