# mail/.config/systemd/user/offlineimap-oneshot.timer.d/override.conf
[Timer]
OnStartupSec=5m
OnUnitInactiveSec=
OnUnitActiveSec=30m
AccuracySec=5m