# media/.config/systemd/user/mpd-notification.service.d/override.conf
[Service]
ExecStart=
ExecStart=/home/jean/.local/bin/mpd-notification

[Unit]
After=dunst.service
Require=dunst.service
